# scripts/wallpaper_set.cron
*/10 * * * *  /home/julian/dotfiles/scripts/wallpaper_changer_cron.sh
